# vagrant/ansible/roles/zamsync/templates/zamsync-serve.service.j2
[Unit]
Description=ZamSync Sync Server ({{ zamsync_role }})
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={{ zamsync_user }}
Group={{ zamsync_group }}
WorkingDirectory={{ zamsync_data_dir }}
ExecStart={{ zamsync_install_dir }}/zamsync serve {{ zamsync_data_dir }} 0.0.0.0:{{ zamsync_sync_port }} --tls --policy {{ zamsync_policy }}
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ zamsync_data_dir }}

[Install]
WantedBy=multi-user.target
